# kueui/root/etc/services.d/app/finish
#!/usr/bin/execlineb -S0

# if supervised process's exit code is neither 0 or 256, bring container down
if { s6-test ${1} -ne 0 }
if { s6-test ${1} -ne 256 }

s6-svscanctl -t /var/run/s6/services
